# Library/LaunchAgents/io.github.vim-denops.LaunchAtLogin.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>io.github.vim-denops.LaunchAtLogin</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-c</string>
      <string>source $HOME/.zsh/.zshenv && [ "$ENABLE_DENOPS_SHARED_SERVER" = "1" ] && deno run -q --no-lock -A $HOME/.local/share/nvim/lazy/denops.vim/denops/@denops-private/cli.ts --hostname 127.0.0.1 --port 32123</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
  </dict>
</plist>
